# root/etc/services.d/discord-bot/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc python3 /app/bot.py -d /config
